# deployment/wallboard.service
[Unit]
Description=WP Community Wallboard (lokal API-proxy + statisk frontend)
Documentation=https://github.com/jesperlowe/Wp-Community-Wallboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=wallboard
Group=wallboard
WorkingDirectory=/opt/wallboard
EnvironmentFile=/opt/wallboard/.env
ExecStart=/usr/bin/env node server/server.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wallboard

# Hærdning — servicen har kun brug for at læse sin egen kode og skrive til
# /var/lib/wallboard (cachen). Alt andet af filsystemet er skrivebeskyttet.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/wallboard
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
