# deployment/wallboard.service
[Unit]
Description=WP Community Wallboard (lokal API-proxy + statisk frontend)
Documentation=https://github.com/jesperlowe/Wp-Community-Wallboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=wallboard
Group=wallboard
WorkingDirectory=/opt/wallboard
EnvironmentFile=/opt/wallboard/.env
ExecStart=/usr/bin/env node server/server.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wallboard

# Hærdning — servicen har kun brug for at læse sin egen kode og skrive til
# /var/lib/wallboard (cachen). Alt andet af filsystemet er skrivebeskyttet.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/wallboard
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true

# [2026-07-16, bugfix] MemoryDenyWriteExecute=true (fjernet) er uforenelig
# med Node.js' V8-motor: JIT-kompilering (Baseline/Sparkplug-compileren)
# allokerer kode-sider som skrivbare, skriver bytecode, og mprotect()'er dem
# derefter til eksekverbare — den flertrins W^X-veksling er præcis hvad
# denne seccomp-hærdning blokerer. Symptom var en total crash-loop ved
# opstart: "Fatal error ... Check failed: 12 == errno" (ENOMEM) inde i
# v8::internal::MemoryAllocator::SetPermissionsOnExecutableMemoryChunk(),
# efterfulgt af SIGTRAP og gentagne systemd-genstarter. Se systemd.exec(5):
# "this option is not suitable for services that generate program code
# dynamically, for example JIT engines" — Node.js er præcis sådan en motor.
# De øvrige hærdningsdirektiver herover rører ikke ved JIT'en og er uændrede.

[Install]
WantedBy=multi-user.target
